# config/includes.chroot/etc/xdg/autostart/huayra-servidor-escolar-conky.desktop
[Desktop Entry]
Type=Application
Name=Informacion en el escritorio
Exec=conky -c /usr/share/huayra-servidor-escolar/conkyrc -p 15
